# install/claude-notify.service
[Unit]
Description=Claude Code notification daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/claude-notify daemon
Restart=on-failure
RestartSec=10
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=HOME=%h
Environment=AWS_PROFILE=your-aws-profile

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/claude-notify
ReadWritePaths=/run/user/%U/claude-notify
ReadWritePaths=%h/.config/claude-notify
ReadWritePaths=%h/.aws
PrivateTmp=true

[Install]
WantedBy=default.target
